# course 2/Daemon/backup-daemon.service
[Unit]
Description=Backup Daemon

[Service]
ExecStart=
Restart=always

[Install]
WantedBy=default.target
