# course 2/Daemon/backup-daemon.service
[Unit]
Description=Backup Daemon

[Service]
ExecStart=/home/nikolay/Backup/build/Backup
Restart=always

[Install]
WantedBy=default.target
